# sensu-docker-client/conf/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[program:sensu-client]
command=/bin/bash -c "exec /opt/sensu/bin/sensu-client -d $CONFIG_DIR -e $EXTENSION_DIR -p $CLIENT_PID -l $LOG_FILE -L $LOG_LEVEL $OPTIONS"

[include]
files = supervisord.d/*.ini
